# proxy/haproxy-node/conf/supervisor/haproxy.conf
# haproxy.
[program:haproxy]
command=/usr/sbin/haproxy -f /etc/haproxy/haproxy.cfg
autostart=true
autorestart=true
stdout_events_enabled=true
stderr_events_enabled=true
