# ci/php5.6/conf/supervisor/ci.conf
# MySQL.
[program:mysql]
command=/usr/bin/pidproxy /var/run/mysqld/mysqld.pid /usr/bin/mysqld_safe --pid-file=/var/run/mysqld/mysqld.pid
autostart=true
autorestart=true
stdout_events_enabled=true
stderr_events_enabled=true

# Disable socat.
[program:socat]
autostart=false
autorestart=false
